# service/com.csend.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd (macOS) : lance `csend serve` (nœud du bus) au login, en TLS hybride PQC.
     Adapte le chemin du binaire si besoin (which csend). Installe-le ainsi :
       cp service/com.csend.serve.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.csend.serve.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.csend.serve</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/csend</string>
    <string>serve</string>
    <string>--tls</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/csend-serve.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/csend-serve.err</string>
</dict>
</plist>
