# service/csend-serve.service
# systemd (Linux) : nœud du bus csend en service utilisateur, TLS hybride PQC.
# Installe-le ainsi :
#   mkdir -p ~/.config/systemd/user
#   cp service/csend-serve.service ~/.config/systemd/user/
#   systemctl --user enable --now csend-serve
[Unit]
Description=csend — bus de messages inter-agents (listener)
After=network.target

[Service]
ExecStart=%h/.local/bin/csend serve --tls
Restart=on-failure
RestartSec=2
# Durcissement (§38)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
